# files/systemd/system/idwt-nu-system.path
[Unit]
Description=Run IDWT (nu) Apply (System) when config changes

[Path]
PathModified=/etc/idwt-nu/config.yml

[Install]
WantedBy=multi-user.target
